# com.user.voice-input.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.voice-input</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/waynesun/voice-input/venv/bin/python</string>
        <string>/Users/waynesun/voice-input/voice_input.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/voice-input.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voice-input.err</string>
</dict>
</plist>
